# crontab
*/5 * * * * pipenv run python -u /home/jobrunner/import.py
